# .config/systemd/user/random_hyprpaper.timer
[Unit]
Description=Set a random wallpaper every minute

[Timer]
Persistent=true
# Every Minute
#OnCalendar=*:0/1

# Every 30 Minutes
OnCalendar=*:0/30

# Every Hour
#OnCalendar=*:00

[Install]
WantedBy=timers.target
