# config/systemd/user/random_hyprpaper.service
[Unit]
Description=Set a random wallpaper with waypaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/random_hyprpaper.sh
